# templates/frp.service.erb
[Unit]
Description=<%= description%>
After=network.target

[Service]
Type=simple
User=nobody
Restart=on-failure
RestartSec=5s
ExecStart=<%= path %> -c <%= @server_config_path %>
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target